# launchd/com.arrival.notify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.arrival.notify</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/gimhyeongjun/Desktop/codex code/scripts/run-notify-once.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/gimhyeongjun/Desktop/codex code</string>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>21</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gimhyeongjun/Desktop/codex code/logs/launchd-notify.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gimhyeongjun/Desktop/codex code/logs/launchd-notify.err.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
